# etc/services.d/spoofy/run
#!/usr/bin/with-contenv sh
yarn --cwd /app start:prod