# etc/services.d/spoofy/run
#!/usr/bin/with-contenv sh
exec yarn --cwd /app start:prod