# etc/services.d/pulse/run
#!/usr/bin/execlineb
pulseaudio --exit-idle-time=-1 --log-level=error